# elasticsearch/files/elasticsearch.service
[Service]
LimitMEMLOCK=infinity

# Restart On Failure
StartLimitInterval=20s
Restart=always
RestartSec=3
